# _docker/allin/supervisord.conf
[supervisord]
nodaemon=true
logfile_maxbytes=50MB
logfile_backups=5

[program:pos]
command=/usr/bin/meter --network %(ENV_NETWORK)s --verbosity 3 --data-dir /pos --api-addr 0.0.0.0:8669 --api-cors * %(ENV_POS_EXTRA)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/var/log/supervisor/pos.log
 
[program:pow]
command=/usr/bin/bitcoind -datadir=/pow -skipibd %(ENV_POW_EXTRA)s
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile = /var/log/supervisor/pow.log

[program:rsyslog]
command=service rsyslog start